# contrib/razerd-battery-notify.service
[Unit]
Description=Notify on low Razer mouse battery

[Service]
Type=oneshot
ExecStart=%h/.local/bin/razerd-battery-notify
